# armbian-build/userpatches/overlay/etc/systemd/system/provision-node.service
[Unit]
Description=First Boot Node Provisioning
After=network-online.target
Wants=network-online.target
# 已跑过则标记文件存在，ConditionPathExists 取反确保不重复执行
ConditionPathExists=!/etc/hive/provisioned

[Service]
Type=oneshot
ExecStart=/usr/local/bin/provision-node.sh
RemainAfterExit=no
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
